# camara-ai-detector.service
[Unit]
Description=Cámara AI Detector de Personas
After=network.target
Wants=network.target
# Requires=camera.target

[Service]
Type=simple
User=root111
WorkingDirectory=/home/root111/camara_1
Environment=PATH=/home/root111/camara_1/venv_camara_ai/bin
ExecStart=/home/root111/camara_1/venv_camara_ai/bin/python3 detector_entrada_salida_v2.py --config config_detector.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Configuración de recursos
LimitNOFILE=65536
LimitNPROC=4096

# Configuración de seguridad simplificada

[Install]
WantedBy=multi-user.target 